# supervisor/docker.ini
[program:dockerd]
command=/bin/bash -c '/usr/local/bin/dockerd ${DOCKER_ARGS}'
autostart=true
autorestart=true
stderr_logfile=/var/log/dockerd.err.log
stdout_logfile=/var/log/dockerd.out.log
